# supervisor/app.conf
[program:webapp]
command=/usr/bin/runapp.sh
autostart=true
autorestart=true
